# services/rtp-receiver-2.0.service
[Service]
Type=simple
Restart=always
RestartSec=5
ExecStart=/usr/bin/gst-launch-1.0 udpsrc port=46000 buffer-size=65536 ! \
    application/x-rtp,media=audio,payload=127,clock-rate=48000,encoding-name=L16,channels=2 ! \
    rtpjitterbuffer latency=200 ! \
    rtpL16depay ! \
    audioconvert ! \
    pulsesink stream-properties="device=extend_2.0_4.1"

[Install]
WantedBy=default.target
